# conf/letsencrypt-renew.cron
#!/bin/sh
#
# certbot renewal cron job
#
# Dockerfile.nginx.letsencrypt copies to /etc/cron.daily/aai-letsencrypt-renew

test -x /usr/bin/certbot || exit 0
/usr/bin/certbot renew -w /opt/starbug.com/www --renew-hook "service nginx reload"
